# dashboard-watchdog.service
[Unit]
Description=Dashboard Watchdog

[Service]
ExecStart=/home/pi/dashboard/watchdog.sh
